# llm_collector/com.kevinharlan.llmcollector.reset.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kevinharlan.llmcollector.reset</string>

  <!-- Run every day at midnight local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>0</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/kevinharlan/llm_collector/reset_collector.sh</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/kevinharlan/llm_collector/reset_launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kevinharlan/llm_collector/reset_launchd.err</string>

  <key>RunAtLoad</key><true/>
</dict>
</plist>
